# autostart/download_sorter_autostart.desktop
[Desktop Entry]
Type=Application
Name=Silent Login Script
Exec=/bin/bash -c "~/.config/automation/download_sorter.sh &"
X-GNOME-Autostart-enabled=true
NoDisplay=true
