# systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --skip-login --nonewline --noissue --autologin oreo --noreset --noclear %I $TERM
Type=simple
TTYVTDisallocate=no
